# harbor/sentry/harbor-sentry.service
[Unit]
Description=Harbor Sentry

[Service]
Type=simple
ExecStart=/usr/local/bin/harbor-sentry

Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
